# ui/polymarket-ui-backend.service
[Unit]
Description=Polymarket Trading Bot - Backend API
After=network.target

[Service]
Type=simple
User=Gabe
WorkingDirectory=/home/Gabe/gubu-workspace/tmp/repos/polymarket-arbitrage-bot/ui/backend
ExecStart=/usr/bin/python app/main.py
Restart=always
RestartSec=10
StandardOutput=append:/home/Gabe/gubu-workspace/tmp/repos/polymarket-arbitrage-bot/ui/backend.log
StandardError=append:/home/Gabe/gubu-workspace/tmp/repos/polymarket-arbitrage-bot/ui/backend.log

[Install]
WantedBy=default.target
